# system_related/network/gateway-network-monitor.service
[Unit]
Description=Zero Axis gateway network monitor
Requires=gateway-startup.service
After=gateway-startup.service gateway-network-apply.service

[Service]
Type=simple
ExecStart=/opt/gateway/scripts/gateway-network-monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
