# system_related/aes/gateway-aes.service
[Unit]
Description=Zero Axis Gateway Analytics Engine
Requires=gateway-startup.service redis-server.service
Wants=gateway-network-apply.service gateway-network-monitor.service gateway-system-metrics.service
After=gateway-startup.service gateway-network-apply.service gateway-network-monitor.service gateway-system-metrics.service gateway-cm5-boot-chirp.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/gateway/aes_bin
Environment=METACRUST_GATEWAY_ROOT=/opt/gateway
Environment=METACRUST_STORAGE_ROOT=/opt/gateway/software_storage
ExecStart=/opt/gateway/aes_bin/gateway-aes
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
